# vps/cloudflared/cloudflared.service
# =============================================================================
# Cloudflare Tunnel — systemd 常駐單元（放 /etc/systemd/system/cloudflared.service）
# 對應 config.yml 的最佳實踐佈署。RunAs 用專用低權限帳號，不跑 root。
# =============================================================================
[Unit]
Description=Cloudflare Tunnel for esggo VPS
After=network.target

[Service]
Type=simple
User=cloudflared
Group=cloudflared
ExecStart=/usr/local/bin/cloudflared tunnel run --config /etc/cloudflared/config.yml esggo-vps
Restart=on-failure
RestartSec=5
# 憑證置於 /etc/cloudflared，權限 600，僅 cloudflared 群組可讀
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/cloudflared

[Install]
WantedBy=multi-user.target
